# systemd-files-podman/garage-webui.service
# container-garage-webui.service
# autogenerated by Podman 5.6.0
# Sun Apr  5 13:01:37 UTC 2026

[Unit]
Description=Podman container-garage-webui.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
TimeoutStopSec=70
ExecStart=/usr/bin/podman run \
        --cidfile=%t/%n.ctr-id \
        --cgroups=no-conmon \
        --rm \
        --sdnotify=conmon \
        -d \
        --replace \
        --name=garage-webui \
        --pod=pod_garage \
        --label io.podman.compose.config-hash=6342d71b626060fa85a7a22b5370699cd89539fd8c7b91e47428f8dde851862c \
        --label io.podman.compose.project=garage \
        --label io.podman.compose.version=1.5.0 \
        --label PODMAN_SYSTEMD_UNIT=podman-compose@garage.service \
        --label com.docker.compose.project=garage \
        --label com.docker.compose.project.working_dir=/home/shopbot/garage \
        --label com.docker.compose.project.config_files=docker-compose.yml \
        --label com.docker.compose.container-number=1 \
        --label io.podman.compose.service=webui \
        --label com.docker.compose.service=webui \
        -e API_BASE_URL=http://garage:4903 \
        -e S3_ENDPOINT_URL=http://garage:4900 \
        -v /home/shopbot/garage/garage.toml:/etc/garage.toml:Z \
        --network=garage_net:alias=webui \
        -p 4909:3909 docker.io/khairul169/garage-webui:latest
ExecStop=/usr/bin/podman stop \
        --ignore -t 10 \
        --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
        -f \
        --ignore -t 10 \
        --cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target