# systemd-files-podman/garage.service
# container-garage.service
# autogenerated by Podman 5.6.0
# Sun Apr  5 13:01:37 UTC 2026

[Unit]
Description=Podman container-garage.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
TimeoutStopSec=70
ExecStart=/usr/bin/podman run \
        --cidfile=%t/%n.ctr-id \
        --cgroups=no-conmon \
        --rm \
        --sdnotify=conmon \
        -d \
        --replace \
        --name=garage \
        --pod=pod_garage \
        --label io.podman.compose.config-hash=6342d71b626060fa85a7a22b5370699cd89539fd8c7b91e47428f8dde851862c \
        --label io.podman.compose.project=garage \
        --label io.podman.compose.version=1.5.0 \
        --label PODMAN_SYSTEMD_UNIT=podman-compose@garage.service \
        --label com.docker.compose.project=garage \
        --label com.docker.compose.project.working_dir=/home/shopbot/garage \
        --label com.docker.compose.project.config_files=docker-compose.yml \
        --label com.docker.compose.container-number=1 \
        --label io.podman.compose.service=garage \
        --label com.docker.compose.service=garage \
        -v /home/shopbot/garage/garage.toml:/etc/garage.toml:Z \
        -v /home/shopbot/garage/meta:/var/lib/garage/meta:Z \
        -v /home/shopbot/garage/data:/var/lib/garage/data:Z \
        --network=garage_net:alias=garage \
        -p 4900:3900 \
        -p 4901:3901 \
        -p 4903:3903 docker.io/dxflrs/garage:v2.2.0
ExecStop=/usr/bin/podman stop \
        --ignore -t 10 \
        --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
        -f \
        --ignore -t 10 \
        --cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target