# systemd-files-podman/postgres.service
# container-postgres.service
# autogenerated by Podman 5.6.0
# Sun Apr  5 13:01:37 UTC 2026

[Unit]
Description=Podman container-postgres.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
TimeoutStopSec=70
ExecStart=/usr/bin/podman run \
        --cidfile=%t/%n.ctr-id \
        --cgroups=no-conmon \
        --rm \
        --sdnotify=conmon \
        -d \
        --replace \
        --name=postgres \
        --pod=pod_postgres \
        --label io.podman.compose.config-hash=d54fab614a56d428290ae681400b7d2aab17dfd5e6d48561fe6672000e5da853 \
        --label io.podman.compose.project=postgres \
        --label io.podman.compose.version=1.5.0 \
        --label PODMAN_SYSTEMD_UNIT=podman-compose@postgres.service \
        --label com.docker.compose.project=postgres \
        --label com.docker.compose.project.working_dir=/home/shopbot/postgres \
        --label com.docker.compose.project.config_files=postgres.yml \
        --label com.docker.compose.container-number=1 \
        --label io.podman.compose.service=postgres \
        --label com.docker.compose.service=postgres \
        -e POSTGRES_USER=postgres \
        -e POSTGRES_PASSWORD=xxx \
        -e POSTGRES_DB=postgres \
        -v /home/shopbot/postgres/data:/var/lib/postgresql:Z \
        --network=pg_net:alias=postgres \
        -p 6432:5432 \
        --shm-size 128mb docker.io/library/postgres:18-alpine
ExecStop=/usr/bin/podman stop \
        --ignore -t 10 \
        --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
        -f \
        --ignore -t 10 \
        --cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target